# huawei-mqtt-publisher/rootfs/etc/service.d/huawei-mqtt/run
#!/bin/sh
cd /
exec python3 huawei_mqtt.py