# huawei-mqtt-publisher/rootfs/etc/service.d/huawei-mqtt/run
#!/usr/bin/with-contenv sh
exec /opt/venv/bin/python -u /app/huawei_mqtt.py