# huawei-mqtt-publisher/rootfs/etc/services.d/huawei-mqtt/run
#!/usr/bin/with-contenv /bin/sh
set -eu

OPTS="/data/options.json"

# Función auxiliar para leer valores del config.json
val() {
  jq -r ".${1} // empty" "$OPTS"
}

# Variables requeridas
export INVERTER_IP="$(val inverter_ip)"
export MQTT_HOST="$(val mqtt_host)"
export MQTT_USERNAME="$(val mqtt_username)"
export MQTT_PASSWORD="$(val mqtt_password)"
export MQTT_PORT="$(val mqtt_port)"
export MODBUS_SLAVE_ID="$(val modbus_slave_id)"
export MODBUS_PORT="$(val modbus_port)"
export MQTT_QOS="$(val mqtt_qos)"
export MQTT_CLIENT_ID="$(val mqtt_client_id)"

# Variables opcionales
PROTO="$(val mqtt_protocol)"
TLS="$(val mqtt_tls)"
KEEPALIVE="$(val mqtt_keepalive)"

[ -n "$PROTO" ] && export MQTT_PROTOCOL="$PROTO"
[ -n "$TLS" ] && export MQTT_TLS="$TLS"
[ -n "$KEEPALIVE" ] && export MQTT_KEEPALIVE="$KEEPALIVE"

echo "[INFO] Starting Huawei Inverter MQTT Publisher..."
echo "[INFO] Inverter: $INVERTER_IP | MQTT: $MQTT_HOST:$MQTT_PORT (client $MQTT_CLIENT_ID)"

exec /opt/venv/bin/python -u /app/huawei_mqtt.py